# files/etc/hotplug.d/net/90-usb-tether
#!/bin/sh

[ "${ACTION}" = "add" ] || exit 0
[ -n "${DEVICENAME}" ] || exit 0
[ -d "/sys/class/net/${DEVICENAME}" ] || exit 0

DRIVER_LINK="/sys/class/net/${DEVICENAME}/device/driver"
[ -L "${DRIVER_LINK}" ] || exit 0

DRIVER_NAME="$(basename "$(readlink -f "${DRIVER_LINK}")")"

case "${DRIVER_NAME}" in
    cdc_ether|cdc_ncm|ipheth|rndis_host)
        ;;
    *)
        exit 0
        ;;
esac

changed=0
current_device="$(uci -q get network.usb_tether.device || true)"
if [ "${current_device}" != "${DEVICENAME}" ]; then
    uci set network.usb_tether.device="${DEVICENAME}"
    changed=1
fi

wan_zone="$(uci show firewall | sed -n "s/^firewall\.\(@zone\[[0-9]\+\]\)\.name='wan'.*/\1/p" | head -n 1)"
if [ -n "${wan_zone}" ]; then
    wan_networks="$(uci -q get "${wan_zone}.network" || true)"
    echo "${wan_networks}" | tr ' ' '\n' | grep -Fx 'usb_tether' >/dev/null 2>&1 || {
        uci add_list "${wan_zone}.network=usb_tether"
        changed=1
    }
fi

if [ "${changed}" -eq 1 ]; then
    uci commit network
    [ -n "${wan_zone}" ] && uci commit firewall
    /etc/init.d/network reload >/dev/null 2>&1 || true
fi

ifup usb_tether >/dev/null 2>&1 || true

exit 0
